# automation/bluetooth-connected.service
[Unit]
Description=Wait for Bluetooth device connection
After=bluetooth.target
Requires=bluetooth.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_bluetooth_connected.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
